# files/root/etc/supervisor/conf.d/tcp443.conf
[program:tcp443]
directory=/etc/openvpn/
command=/usr/sbin/openvpn --config tcp443.conf
environment=USER="root",TERM="dumb"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
